# openosint-mcp.service
[Unit]
Description=OpenOSINT MCP Server Service
Documentation=https://openosint.github.io/OpenOSINT
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/openosint
EnvironmentFile=/srv/openosint/.env
ExecStart=/usr/bin/python3.12 -m openosint.mcp_server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/srv/openosint/reports /srv/openosint/openosint.db /srv/openosint/.cache
PrivateTmp=true
ProtectHome=read-only

[Install]
WantedBy=default.target
